# config/systemd/monit-intel-ingest.service
[Unit]
Description=Monit-Intel Ingest Service
Documentation=https://github.com/heverz/monit-intel
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=heverz
WorkingDirectory=/home/heverz/py_projects/monit-intel

# Set PYTHONPATH for module discovery
Environment="PYTHONPATH=/home/heverz/py_projects/monit-intel/src"

# Load environment from .env (optional - can also be set in monit-intel-agent.service)
EnvironmentFile=-/home/heverz/.env

# Run ingest once (timer handles 5-minute scheduling)
ExecStart=/home/heverz/.pixi/bin/pixi run python -m monit_intel.ingest

# Restart on failure
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=monit-intel-ingest

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/home/heverz/py_projects/monit-intel

[Install]
WantedBy=multi-user.target
